# systemd.service
# /etc/systemd/system/mlendpoints_bun.service
# SvelteKit Project Service (Bun via Unix Socket)
#
# Architecture: Caddy → Unix Socket → Bun → SvelteKit
#
# Commands:
#   sudo systemctl status mlendpoints_bun
#   sudo journalctl -u mlendpoints_bun -f
#   sudo systemctl restart mlendpoints_bun

[Unit]
Description=MLEndpoints SvelteKit Project (Bun via Unix Socket)
After=network.target

[Service]
Type=simple

User=ubuntu
Group=www-data
WorkingDirectory=/var/www/ml-endpoints

# Loads SOCKET_PATH and other config from .env
EnvironmentFile=/var/www/ml-endpoints/.env
Environment="NODE_ENV=production"
# Explicitly set SOCKET_PATH to ensure it's available (systemd EnvironmentFile can be unreliable)
Environment="SOCKET_PATH=/run/bun/mlendpoints_bun.sock"

# Creates /run/bun/ owned by ubuntu, cleaned on stop
RuntimeDirectory=bun
RuntimeDirectoryMode=0750

ExecStart=/home/ubuntu/.bun/bin/bun run build/index.js

# Restart on crash, 3s delay
Restart=always
RestartSec=3

# Memory limits for SvelteKit
MemoryHigh=512M
MemoryMax=1G

# Basic security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/www/ml-endpoints /run/bun
PrivateTmp=true

[Install]
WantedBy=multi-user.target
